# docker/crontab
* * * * * www-data cd /var/www/html && php email_worker.php >> /var/www/html/storage/logs/email_worker.log 2>&1
*/2 * * * * www-data cd /var/www/html && php pdf_worker.php >> /var/www/html/storage/logs/pdf_worker.log 2>&1
*/15 * * * * www-data cd /var/www/html && php payout_worker.php >> /var/www/html/storage/logs/payout_worker.log 2>&1
